# deploy/systemd/pica.service
# Pica — systemd service (Linux)
#
# Install to /etc/systemd/system/pica.service, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pica
#
# Runs Pica as an unprivileged user. The passphrase is read from an
# EnvironmentFile so the service can restart unattended — an interactive boot
# would block on the "Passphrase:" prompt. Full walkthrough: docs/deployment.md.

[Unit]
Description=Pica time management
After=network.target

[Service]
Type=simple
User=pica
Group=pica
WorkingDirectory=/opt/pica
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Passphrase lives here, NOT in this unit file. Create it root-owned, mode 600:
#   echo 'PICA_PASSPHRASE=your-passphrase' | sudo tee /etc/pica.env
#   sudo chmod 600 /etc/pica.env
EnvironmentFile=/etc/pica.env

# Hardening — safe with Pica's file storage as long as data/ and backups/
# live under WorkingDirectory (which stays writable).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
